# src/systemd/sd-inhibit-bridge.service.in
[Unit]
Description=An idle inhibition daemon that forwards inhibitor locks to logind

[Service]
Type=simple
ExecStart=@SDIB_BIN@
Restart=on-failure
Slice=session.slice
StandardError=journal

[Install]
WantedBy=default.target
